# conf/crontab
# Backup DB
0 0 * * * /srv/http/www/django-nuxt-template/manage.sh backupdb